# docker/config/conf.d/cron.d/dmarc
#
#  dmarc_reports
#
# m h dom mon dow user  command
*/5 * * * * root dmarc_reports.sh
