# docker/config/conf.d/cron.d/opendmarc_reports
#
#  dmarc_reports
#
# m h dom mon dow user  command
*/5 * * * * root /usr/local/bin/dmarc_reports.sh

