# docker/conf/cron.d/freshclam
#
#  Freshclam
#
# m h dom mon dow user  command
  0 5  21  *   *  root  freshclam-update.sh

